# debs/zabbix-check/zabbix-check-confixx/etc/cron.d/zabbix-check-confixx
*/10 * * * * root /root/confixx/confixx_counterscript.pl 2>/dev/null; echo $? > /usr/share/zabbix/confixx/lizenz
