# debs/zabbix-check/zabbix-check-dpkg/etc/cron.d/zabbix-check-dpkg
#zabbix-check-dpkg
35 */2 * * * root perl -e 'sleep rand(1800)' && ionice -c3 nice -19 /usr/bin/zabbix-check-dpkg.sh
