# debs/zabbix-check/zabbix-check-dpkg/etc/cron.d/zabbix-check-dpkg
#zabbix-check-dpkg
35 * * * * root ionice -c3 nice -19 /usr/bin/zabbix-check-dpkg.sh cron
