# externalscripts/zabbix-check-ssllabs/cron.d/zabbix-check-ssllabs
# Renames every host found in /etc/zabbix/zabbix-check-ssllabs.db
0 1     * * *    root   /etc/zabbix/externalscripts/zabbix-check-ssllabs.lua -u