# debs/zabbix-agent-config/debian/postinst
#!/bin/sh
# postinst script for cubos-fw
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in

  configure)

	if [ ! -f /etc/zabbix/zabbix_agentd.conf ]; then
		echo "ERROR: NO zabbix_agentd.conf CONFIG"
		exit 2
	fi

	sed -i 's/^Server=.*/Server=zabbix.caix.de, 213.131.245.67, 10.192.0.30/g' /etc/zabbix/zabbix_agentd.conf
	sed -i 's/^Hostname=.*/Hostname='`cat /etc/hostname`'/g' /etc/zabbix/zabbix_agentd.conf

	if grep '^Include' /etc/zabbix/zabbix_agentd.conf >/dev/null; then
		sed -i 's%^Include=.*%Include=/etc/zabbix/zabbix_agentd.d/%g' /etc/zabbix/zabbix_agentd.conf
	else
		echo "Include=/etc/zabbix/zabbix_agentd.d/" >> /etc/zabbix/zabbix_agentd.conf
	fi

	#fix for old systems
	usermod -s /bin/bash zabbix

	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d zabbix-agent restart || true
	else
		/etc/init.d/zabbix-agent restart || true
	fi	

  ;;

  abort-upgrade|abort-remove|abort-deconfigure)
    exit 0
  ;;

  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
  ;;

esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

