# debs/zabbix-check/zabbix-check-postfix/etc/cron.d/zabbix-check-postfix
#zabbix-check-postfix logger
*/5 * * * * root ionice -c3 nice -19 /usr/lib/zabbix-check/bin/zabbix-check-postfix.sh
