# debs/zabbix-check/zabbix-check-smart/debian/postinst
#!/bin/sh -e

case "$1" in
    configure)
        chmod -v +x /usr/sbin/zabbix_check_smart_get.py
        chmod -v +x /usr/sbin/zabbix_check_smart_discover.py
        chmod -v 440 /etc/sudoers.d/zabbix-check-smart

        #remove old version of sudo rights (pre 1:0.4)
        if grep -q 'zabbix_check_smart_get.py' /etc/sudoers; then
            grep -v 'zabbix_check_smart_get.py' /etc/sudoers > /etc/.sudoers
            mv /etc/.sudoers /etc/sudoers
            chmod 440 /etc/sudoers
        fi

        echo "Enable SMART if disabled"
        for hdd in $(ls /dev/ | egrep ^[hs]{1}d[a-z]{1}$); do
            if [ "$(smartctl -i /dev/$hdd | grep Disabled | wc -l)" -gt "0" ]; then
                smartctl -s on -q silent /dev/$hdd
                echo "Enabled SMART on $hdd"
            fi
        done

        if which invoke-rc.d >/dev/null 2>&1; then
            invoke-rc.d zabbix-agent stop || true
            invoke-rc.d zabbix-agent start || true
        else
            /etc/init.d/zabbix-agent stop || true
            /etc/init.d/zabbix-agent start || true
        fi
        ;;

    abort-upgrade|abort-deconfigure|abort-remove)
	    ;;

    *)
	    echo "$0 called with unknown argument \`$1'" 1>&2
	    exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
