# debs/zabbix-check/zabbix-check-fs/etc/cron.d/zabbix-check-fsck
#zabbix-check-fs
30 */6 * * * root /usr/lib/zabbix-check/bin/zabbix-check-fsck.sh > /usr/share/zabbix/fsck
