# debs/zabbix-check/zabbix-check-fs/etc/cron.d/zabbix-check-fsck
#zabbix-check-fs logger
30 */6 * * * root /usr/bin/zabbix-check-fs.sh > /usr/share/zabbix/fsck
